# ntripcaster.service
[Unit]
Description=NTRIP Caster Service
After=network.target etcd.service rqlite.service
Documentation=https://github.com/your-org/ntripcaster

[Service]
User=ntrip
Group=ntrip
WorkingDirectory=/opt/ntripcaster
Environment=RUST_LOG=info
Environment=CONFIG_PATH=/etc/ntripcaster/config.toml
ExecStart=/opt/ntripcaster/bin/ntripcaster
Restart=always
RestartSec=5s
LimitNOFILE=65535

# 资源限制
MemoryLimit=2G
CPUQuota=80%

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/var/log/ntripcaster /var/lib/ntripcaster

[Install]
WantedBy=multi-user.target